# cliproxyapi/rootfs/etc/services.d/cliproxyapi/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting CLIProxyAPI server on :8317"

exec /usr/bin/CLIProxyAPI --config /config/cliproxyapi/config.yaml
